# archiso/airootfs/etc/xdg/autostart/streekx.desktop
[Desktop Entry]
Type=Application
Name=Streekx World Launcher
Exec=/usr/bin/streekx-start.sh
Icon=system-run
X-GNOME-Autostart-enabled=true
